# packages/ca-server/debian/ca-server.ca-server-init.service
[Unit]
Description=Initialize CA (Easy-RSA) on first boot
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/ca-server-init
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

