# packages/backup-system-config/linux-final-backup.timer
[Unit]
Description=Run Linux Final backup daily

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=5m
Unit=linux-final-backup.service

[Install]
WantedBy=timers.target
