# packages/backup-system-config/debian/prerm
#!/usr/bin/env bash
set -euo pipefail

if command -v systemctl >/dev/null 2>&1; then
  systemctl disable --now linux-final-backup.timer 2>/dev/null || true
  systemctl disable --now linux-final-backup-rotate.timer 2>/dev/null || true
fi

exit 0
